# docker/ha-addon-rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: ESPHome
# Runs the NGINX proxy
# ==============================================================================

bashio::log.info "Waiting for dashboard to come up..."

while [[ ! -S /var/run/esphome.sock ]]; do
  sleep 0.5
done

bashio::log.info "Starting NGINX..."
exec nginx
