# docker/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: ESPHome
# Runs the NGINX proxy
# ==============================================================================

bashio::log.info "Starting NGINX..."
exec nginx -g "daemon off;"
