# docker/rootfs/etc/services.d/esphome/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: ESPHome
# Runs the ESPHome dashboard
# ==============================================================================

export ESPHOME_IS_HASSIO=true

if bashio::config.true 'leave_front_door_open'; then
    export DISABLE_HA_AUTHENTICATION=true
fi

if bashio::config.true 'streamer_mode'; then
    export ESPHOME_STREAMER_MODE=true
fi

if bashio::config.true 'status_use_ping'; then
    export ESPHOME_DASHBOARD_USE_PING=true
fi

if bashio::config.has_value 'relative_url'; then
    export ESPHOME_DASHBOARD_RELATIVE_URL=$(bashio::config 'relative_url')
fi

bashio::log.info "Starting ESPHome dashboard..."
exec esphome /config/esphome dashboard --socket /var/run/esphome.sock --hassio
